# scripts/systemd/bts-lineup-collect.service
[Unit]
Description=BTS lineup time collection (one-shot)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/bts
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="UV_CACHE_DIR=/tmp/uv-cache"
ExecStart=%h/.local/bin/uv run bts data collect-lineup-times
StandardOutput=journal
StandardError=journal
